# display/80-hdmi.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{XAUTHORITY}="/home/pollo/.Xauthority" RUN+="/usr/bin/hotplug-hdmi"
